# ansible/roles/media_ingest/templates/media-ingest-reconcile.timer.j2
# Managed by Ansible — roles/media_ingest. Do not edit by hand.
[Unit]
Description=Run the media-ingest reconcile sweep periodically

[Timer]
# First run shortly after boot, then every interval. Persistent=true catches up
# one missed run if the box was off at the scheduled time.
OnBootSec=5min
OnUnitActiveSec={{ media_ingest_reconcile_interval }}
Persistent=true

[Install]
WantedBy=timers.target
