# ansible/roles/media_ingest/templates/media-ingest-collector.service.j2
# Managed by Ansible — roles/media_ingest. Do not edit by hand.
[Unit]
Description=upstream source gallery collector (the source protocol source client -> gallery API push)
After=network-online.target {{ media_ingest_mount_unit.stdout }}
Wants=network-online.target
Requires={{ media_ingest_mount_unit.stdout }}
# StartLimitIntervalSec belongs in [Unit] (a [Service] placement is ignored).
StartLimitIntervalSec=0

[Service]
Type=simple
User={{ media_ingest_user }}
Group={{ media_ingest_user }}
EnvironmentFile={{ media_ingest_dir }}/collector.env
ExecStart={{ media_ingest_dir }}/venv/bin/python {{ media_ingest_dir }}/collector.py
Restart=always
RestartSec=5
# the client reconnects on its own; if the process dies, bring it back fast.

# Hardening — the collector only needs its own state/log dirs.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ media_ingest_state }} {{ media_ingest_log_dir }}

[Install]
WantedBy=multi-user.target
