# ansible/roles/media_gallery/templates/collector.service.j2
# Managed by Ansible — roles/media_gallery. Do not edit by hand.
[Unit]
Description=upstream source ephemeral media collector (the source protocol source client listener)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ media_gallery_user }}
Group={{ media_gallery_user }}
EnvironmentFile={{ media_gallery_dir }}/collector.env
ExecStart={{ media_gallery_dir }}/venv/bin/python {{ media_gallery_dir }}/collector.py
Restart=always
RestartSec=5
# source reconnects on its own; if the process dies, bring it back fast.
StartLimitIntervalSec=0

# Hardening — the source client only needs its own state/log dirs.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ media_gallery_state }} {{ media_gallery_log_dir }}

[Install]
WantedBy=multi-user.target
