# ansible/roles/media_gallery/templates/sync.service.j2
# Managed by Ansible — roles/media_gallery. Do not edit by hand.
[Unit]
Description=media-gallery: encrypt+upload captures to Google Drive, prune local cache
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{ media_gallery_user }}
Group={{ media_gallery_user }}
Environment=TG_SAVE_DIR={{ media_gallery_state }}/captures
Environment=TG_RCLONE_REMOTE={{ media_gallery_rclone_remote }}
Environment=TG_CACHE_DAYS={{ media_gallery_cache_days }}
Environment=RCLONE_CONFIG=/home/{{ media_gallery_user }}/.config/rclone/rclone.conf
Environment=INGEST_LOG={{ media_gallery_log_dir }}/sync.log
ExecStart={{ media_gallery_dir }}/sync.sh
