# ansible/roles/media_gallery/templates/upload.service.j2
# Managed by Ansible — roles/media_gallery. Do not edit by hand.
[Unit]
Description=media-gallery upload service (add pics + create folders)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ media_gallery_user }}
Group={{ media_gallery_user }}
Environment=RCLONE_CONFIG=/home/{{ media_gallery_user }}/.config/rclone/rclone.conf
Environment=TG_RCLONE_REMOTE={{ media_gallery_rclone_remote }}
Environment=UPLOAD_PORT={{ media_gallery_upload_port }}
Environment=UPLOAD_BIND={{ ip_media_gallery }}
Environment=TG_DATEMAP_CACHE={{ media_gallery_state }}/datemap.json
Environment=TG_EXCLUDE_FILE={{ media_gallery_state }}/excluded.json
Environment=UPLOAD_MAX_MB={{ media_gallery_upload_max_mb }}
Environment=UPLOAD_STAGING={{ media_gallery_state }}/upload_staging
Environment=REBUILD_DEBOUNCE_SEC=20
# Async GDrive push: how often the background pusher flushes PENDING/ -> GDrive.
Environment=PUSH_INTERVAL_SEC={{ media_gallery_push_interval_sec | default(5) }}
Environment=INGEST_JWKS_URL={{ authentik_url }}/application/o/media-ingest-collector/jwks/
Environment=INGEST_ISSUER={{ authentik_url }}/application/o/media-ingest-collector/
Environment=INGEST_AUDIENCE=media-ingest-collector
ExecStart={{ media_gallery_dir }}/venv/bin/python {{ media_gallery_dir }}/upload_service.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
