# ansible/roles/media_ingest_02/templates/staging.mount.j2
# Managed by Ansible — roles/media_ingest_02. Do not edit by hand.
# tmpfs mount for decrypted-media staging. The staging dir holds DECRYPTED media
# in flight to the gallery. Back it with tmpfs so nothing viewable ever persists
# to disk — a stolen disk / Proxmox backup / snapshot contains no plaintext.
# Unit filename is a systemd-escaped path (computed in tasks/main.yml).
[Mount]
What=tmpfs
Where={{ media_ingest_02_state }}/staging
Type=tmpfs
Options=mode=0750,uid={{ media_ingest_02_user_uid }},gid={{ media_ingest_02_user_gid }},size={{ media_ingest_02_staging_size }}

[Install]
WantedBy=multi-user.target
